# backend/hqmx-backend.service
[Unit]
Description=HQMX Backend - Flask + yt-dlp Service
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/hqmx/backend
Environment="PATH=/home/ubuntu/hqmx/backend/venv/bin"
ExecStart=/home/ubuntu/hqmx/backend/venv/bin/python app.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hqmx-backend

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
